# tools/rc.local
sudo  python /home/pi/bouncer/main.py noserial &